# network.service
[Unit]
Description=Initialize network
DefaultDependencies=no
After=syslog.target udev-settle.service
Wants=network.target udev.service udev-settle.service
Wants=local-fs.target systemd-tmpfiles-setup.service
Before=network.target

[Service]
Type=oneshot
ExecStart=/etc/rc.d/init.d/network start
ExecStop=/etc/rc.d/init.d/network stop
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
