# udev-alsa.rules
SUBSYSTEM=="sound", KERNEL=="controlC?", RUN+="/sbin/alsactl -F restore %n"
