# udev-fake.service
[Unit]
Description=non-udev Coldplug all Devices
DefaultDependencies=no

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/root/fetch_devices
